# 1.13.7-alpine/app.ini
[program:app]
directory=/go
command=/go/app
autostart=true
autorestart=true
stderr_logfile=/go/log/app.err
stdout_logfile=/var/log/app.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0