# 1.13.7-alpine/app.ini
[program:app]
command=/go/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord.log
stdout_logfile_maxbytes=0
logfile_maxbytes=5MB
logfile_backups=100
loglevel=%(LOG_LEVEL)
pidfile=/var/supervisord.pid